# alpine/etc/sv/tinydns/run
#!/usr/bin/env sh
set -eu

exec 2>&1

DAEMON=/usr/bin/tinydns

/sbin/chpst -U tinydns:tinydns -e ./env ${DAEMON}
