# debian/etc/sv/skeleton/run
#!/usr/bin/env sh
set -eu

exec 2>&1

COMMAND=daemon
USER=root

exec /usr/bin/chpst -u ${USER} ${COMMAND}